# systemd/verify_updates.service
[Unit]
Description=Verificar actualizaciones con yay

[Service]
Type=oneshot
ExecStart=/home/kry_/.config/scripts/verify_updates.sh
